# internal/detection/rules/windows_registry_run_key.yaml
id: "windows_registry_run_key"
name: "Registry Run Key Persistence"
description: >
  Detects modification of common autorun registry keys used for persistence.
  Attackers add entries to HKCU\Run or HKLM\Run to survive reboots.
severity: "medium"
type: "threshold"
threshold: 1
window_sec: 60
dedup_window_sec: 3600
mitre_tactics:
  - "TA0003"
mitre_techniques:
  - "T1547.001"
conditions:
  EventType: "windows_registry_set"
  TargetObject|contains|any:
    - "\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\"
    - "\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\"
    - "\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon"
    - "\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Windows\\AppInit_DLLs"
logsource:
  product: windows
  service: sysmon
  category: registry_set
