# internal/detection/rules/windows_scheduled_task_lateral.yaml
id: "windows_scheduled_task_lateral"
name: "Remote Scheduled Task Creation (Lateral Movement)"
description: >
  Detects scheduled tasks created remotely (via schtasks /s or Task Scheduler
  RPC on a remote host) — a common lateral movement technique used by
  ransomware groups and APTs after obtaining credentials.
severity: "high"
type: "threshold"
threshold: 1
window_sec: 60
dedup_window_sec: 1800
mitre_tactics:
  - "TA0008"
  - "TA0002"
mitre_techniques:
  - "T1053.005"
  - "T1021.003"
conditions:
  EventType: "windows_process_create"
  CommandLine|contains|any:
    - "schtasks /create /s "
    - "schtasks.exe /create /s "
    - "Register-ScheduledTask"
  CommandLine|contains: "/s "
logsource:
  product: windows
  service: sysmon
  category: process_creation
